# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/run/supervisord.pid

; ─────────────────────────────────────────────────────────────────────────────
; Nginx — Reverse proxy / static asset server
; ─────────────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ─────────────────────────────────────────────────────────────────────────────
; PHP-FPM — PHP Fast Process Manager
; ─────────────────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=true
startretries=5
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

; ─────────────────────────────────────────────────────────────────────────────
; Laravel Queue Worker
;
; Runs as a second process inside the same container.
; Per SAD ADR-007: for production, this should move to a dedicated Cloud Run
; worker revision. Keep it here for local Docker convenience.
;
; Uses --max-time=3600 so the process is gracefully recycled every hour,
; preventing memory creep in a long-lived process.
; ─────────────────────────────────────────────────────────────────────────────
[program:laravel-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
startretries=3
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15
user=www-data
